# etc/gunicorn/gunicorn.socket
[Unit]
Description=Gunicorn Socket

[Socket]
ListenStream=/var/www/songhill/backend/gunicorn.sock

[Install]
WantedBy=sockets.target

